# rclone_lyo.timer
[Unit]
Description=rclone periodically

[Timer]
Unit=rclone_lyo.service
OnCalendar=*-*-* *:00:00

[Install]
WantedBy=timers.target
